# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:qsshd]
command=/usr/qssh/sbin/sshd -ddd -D